# recipes-core/strbo-recover-system/files/strbo-recovery.service
[Unit]
Description=Reinstall main system from image files

[Service]
Type=oneshot
ExecStart=/usr/libexec/recover-system
StandardOutput=syslog

[Install]
WantedBy=basic.target
