# recipes-multimedia/gerbera/usb-storage-manual/rootfs/etc/systemd/system/mounta.service.d/override.conf
[Service]
ExecStartPre=-/usr/bin/sudo /usr/bin/usb-storage unbind --disable-autoprobe --wait 10 --from /etc/usb-storage-manual
